# services/genieacs/genieacs-services-1.2.8/supervisord.conf
[supervisord]
nodaemon=true
user=genieacs
logfile=/var/log/genieacs/supervisord.log

[program:genieacs-cwmp]
directory=/opt/genieacs
command=/opt/genieacs/dist/bin/genieacs-cwmp
environment=MONGODB_CONNECTION_URL="%(ENV_GENIEACS_MONGODB_CONNECTION_URL)s",REDIS_HOST="%(ENV_GENIEACS_REDIS_HOST)s"
stdout_logfile=/var/log/genieacs/genieacs-cwmp.log
stderr_logfile=/var/log/genieacs/genieacs-cwmp.log
autorestart=true

[program:genieacs-nbi]
directory=/opt/genieacs
command=/opt/genieacs/dist/bin/genieacs-nbi
environment=MONGODB_CONNECTION_URL="%(ENV_GENIEACS_MONGODB_CONNECTION_URL)s",REDIS_HOST="%(ENV_GENIEACS_REDIS_HOST)s"
stdout_logfile=/var/log/genieacs/genieacs-nbi.log
stderr_logfile=/var/log/genieacs/genieacs-nbi.log
autorestart=true

[program:genieacs-fs]
directory=/opt/genieacs
command=/opt/genieacs/dist/bin/genieacs-fs
environment=MONGODB_CONNECTION_URL="%(ENV_GENIEACS_MONGODB_CONNECTION_URL)s",REDIS_HOST="%(ENV_GENIEACS_REDIS_HOST)s"
stdout_logfile=/var/log/genieacs/genieacs-fs.log
stderr_logfile=/var/log/genieacs/genieacs-fs.log
autorestart=true

[program:genieacs-ui]
directory=/opt/genieacs
command=/opt/genieacs/dist/bin/genieacs-ui
environment=MONGODB_CONNECTION_URL="%(ENV_GENIEACS_MONGODB_CONNECTION_URL)s",REDIS_HOST="%(ENV_GENIEACS_REDIS_HOST)s",UI_JWT_SECRET="%(ENV_GENIEACS_UI_JWT_SECRET)s"
stdout_logfile=/var/log/genieacs/genieacs-ui.log
stderr_logfile=/var/log/genieacs/genieacs-ui.log
autorestart=true